# etc/systemd/clean-tmp.service
[Unit]
Description=Delete everything in ~/tmp

[Service]
ExecStart=clean-tmp
User=hamish

[Install]
WantedBy=multi-user.target